# contrib/launchd/com.inzerator.bulk-dmpok.start.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Starts bulk_dmpok.py every night at 22:00. The script holds a lock dir
  under /Volumes/Elements/cuzk-bulk/.lock so a duplicate start (e.g.
  the previous night's process wasn't killed) becomes a no-op + exit.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.inzerator.bulk-dmpok.start</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jan/projekty/inzerator/bulk_dmpok.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jan/projekty/inzerator</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BULK_OUT_DIR</key>
    <string>/Volumes/Elements/cuzk-bulk</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <!--
    Logs go to ~/Library/Logs/ (which launchd user agents can always
    write), NOT /Volumes/Elements — user agents lack Full Disk Access
    to external volumes by default and silently fail with exit 78.
    Note: the script's own per-sheet log still lands on Elements
    (sheets.json + download.log) via BULK_OUT_DIR.
  -->
  <key>StandardOutPath</key>
  <string>/Users/jan/Library/Logs/inzerator/bulk-dmpok-start.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jan/Library/Logs/inzerator/bulk-dmpok-start.err.log</string>
</dict>
</plist>
